# etc/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:statsd]
directory=/statsd
command=/usr/local/node/bin/node /statsd/stats.js /statsd/config.js
autostart=yes
autorestart=yes